# filesystem/etc/supervisor.d/postfix.ini
[program:postfix]
command=/usr/sbin/postfix -c /etc/postfix start
directory=/etc/postfix
process_name=master
startsecs=0
autostart=false
autorestart=false
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
stdout_events_enabled = true
stderr_events_enabled = true
